# deploy/oracle/systemd/expense-ai.service
[Unit]
Description=Expense AI Service (FastAPI + PyTorch)
After=network.target

[Service]
Type=simple
User=%APP_USER%
Group=%APP_USER%
WorkingDirectory=%APP_DIR%/ai_service
EnvironmentFile=/etc/expense-manager/production.env
Environment=PATH=%APP_DIR%/ai_service/.venv/bin
ExecStart=%APP_DIR%/ai_service/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=8
# PyTorch load model có thể mất 30-60s
TimeoutStartSec=180

[Install]
WantedBy=multi-user.target
